# genon/serving/paddle/config/supervisord.conf
[program:paddlex-ocr]
directory=/app
command=/app/.venv/bin/paddlex --serve --pipeline /app/config/ocr.yaml --host 0.0.0.0 --port 8080
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0