# genon/serving/paddle/config/supervisord.conf
[supervisord]
nodaemon=true

[program:paddlex-ocr]
directory=/app
command=/app/.venv/bin/paddlex --serve --pipeline /app/config/ocr.yaml --host 0.0.0.0 --port 8080
autostart=true
autorestart=true
stdout_logfile=/var/log/paddlex/stdout.log
stderr_logfile=/var/log/paddlex/stderr.log